# deploy/supervisor-crm-worker.conf
# /etc/supervisor/conf.d/crm-worker.conf — copy there, then:
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl status crm-worker:*
#
# This is the fix for the queue worker dying silently with nothing to
# restart it (observed repeatedly running it by hand in dev). Supervisor
# keeps it running and restarts it automatically on crash or reboot.

[program:crm-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/crm/backend/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/crm/backend
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=3600
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/crm/backend/storage/logs/worker.log
stopasgroup=true
killasgroup=true

# numprocs=1 is fine for 10 users. If imports/automations start queuing up
# and lagging, bump numprocs (e.g. to 2-3) to run multiple worker processes
# in parallel — Supervisor names them crm-worker_00, crm-worker_01, etc.
#
# --max-time=3600 recycles each worker process hourly (a clean restart,
# not a crash) — guards against slow memory growth in long-running PHP
# workers; Supervisor immediately starts a fresh one, no gap in processing.
